# installer/debian/lib/systemd/system/ftp-simulator-app.service
[Unit]
Description=FTP Simulator — Application server
After=network.target ftp-simulator-db.service
Requires=ftp-simulator-db.service

[Service]
Type=simple
User=ftp-simulator
Group=ftp-simulator
WorkingDirectory=/var/lib/ftp-simulator
EnvironmentFile=-/etc/ftp-simulator/environment
ExecStart=/usr/bin/ftp-backend
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/var/lib/ftp-simulator /var/log/ftp-simulator

[Install]
WantedBy=multi-user.target
