# installer/debian/lib/systemd/system/ftp-simulator-db.service
[Unit]
Description=FTP Simulator — PostgreSQL (Docker)
Documentation=https://github.com/ce-teuf/ftp-simulator
After=docker.service network-online.target
Requires=docker.service
Before=ftp-simulator-app.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/etc/ftp-simulator

# Pull image if not present, then start the compose stack
ExecStartPre=/usr/bin/docker compose -f /etc/ftp-simulator/docker-compose.prod.yml pull --quiet
ExecStart=/usr/bin/docker compose -f /etc/ftp-simulator/docker-compose.prod.yml up -d --wait
ExecStop=/usr/bin/docker compose -f /etc/ftp-simulator/docker-compose.prod.yml stop
ExecReload=/usr/bin/docker compose -f /etc/ftp-simulator/docker-compose.prod.yml restart

Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ftp-simulator-db

[Install]
WantedBy=multi-user.target
